# docker/cron/crontab
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
HOME=/

0 6 * * * root /root/backup.sh > /dev/stdout
